# boot/configure-clocksource.service
[Unit]
Description=Configure kernel clocksource

[Service]
ExecStart=/usr/bin/configure-clocksource

[Install]
WantedBy=multi-user.target
