# templates/al2/runtime/accelerator/nvidia-eula.service
[Unit]
Description=Display Nvidia driver EULA

[Service]
Type=oneshot
ExecStart=/etc/eks/nvidia-eula.sh

[Install]
WantedBy=basic.target