# lessonly-files/aide-startup.service
[Unit]
Description = Check AIDE at startup

[Service]
ExecStart=/bin/aide-startup.sh

[Install]
WantedBy = multi-user.target