# templates/al2023/runtime/rootfs/etc/systemd/system/containerd.service.d/00-runtime-slice.conf
# Must match runtime-cgroups flag set on kubelet for cadvisor metrics to work.
[Service]
Slice=runtime.slice
