# nodeadm/test/e2e/infra/systemd/imds-mock.service
[Unit]
Description=Run the mock IMDS server

[Service]
ExecStart=imds-mock

[Install]
WantedBy=multi-user.target
