# files/chronyd-sync.timer
[Unit]
Description=Half hour chronyd timer

[Timer]
OnCalendar=*:0/30
RemainAfterElapse=yes

[Install]
WantedBy=timers.target
