# files/chronyd-sync.service
[Unit]
Description=Keep system time synced

[Service]
Type=oneshot
ExecStart=/usr/sbin/chronyd -q
