# files/kube-container-runtime-monitor.service
[Unit]
Description=Kubernetes health monitoring for container runtime
After=docker.service

[Service]
Restart=always
RestartSec=10
RemainAfterExit=yes
RemainAfterExit=yes
ExecStartPre=/bin/chmod 544 /usr/local/bin/health-monitor.sh
ExecStart=/usr/local/bin/health-monitor.sh container-runtime

[Install]
WantedBy=multi-user.target
