# files/ntpdate-sync.timer
[Unit]
Description=Hourly ntpdate timer

[Timer]
OnCalendar=hourly
RemainAfterElapse=yes

[Install]
WnatedBy=timers.target
